# start-at-boot.cron
@reboot ( sleep 20; /root/tablet-internet-switch/switch-listener.sh; )
